# containers/webserver/gunicorn.ini
[program:nginx]
command=/usr/sbin/nginx

[program:gunicorn]
directory=/opt/webserver
command=gunicorn -c /opt/gunicorn.conf webservice:app
user=root
autorstart=true
autorestart=true
stdout_logfile=/var/log/webservice.log
stderr_logfile=/var/log/webservice_error.log